# config/supervisord.conf
[supervisord]
nodaemon=true
user=ryot

[program:backend]
command=ryot
environment=

[program:frontend]
command=npx remix-serve ./build/server/index.js
environment=API_URL=http://localhost:5000

[program:reverse-proxy]
command=caddy run --config /etc/caddy/Caddyfile
environment=
